# templates/fairdata-download-server.socket.j2
[Unit]
Description=UNIX socket for Fairdata Download Server

[Socket]
ListenStream=/run/fairdata-download-server.sock
SocketUser=root

[Install]
WantedBy=sockets.target
